# CommApp/CommAppYocto/yocto/poky/meta-myLayer/recipes-iot/iot-client/files/iot-client.service
[Unit]
Description=IoT Client — Edges For Training
Documentation=https://github.com/edgesfortraining/iot-client
After=network-online.target
Wants=network-online.target
# Restart automatically if the server is not yet reachable
StartLimitIntervalSec=0

[Service]
Type=simple
User=root
ExecStart=/usr/bin/iot_client
Restart=on-failure
RestartSec=5s

# Redirect stdout/stderr to journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=iot-client

# Give the process access to GPIO and network
AmbientCapabilities=CAP_NET_RAW
ReadWritePaths=/sys/class/gpio

[Install]
WantedBy=multi-user.target
